# cron/diffsnap.service
[Unit]
Description=Create ZFS snapshots using diffsnap

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/diffsnap
